# ci/production/aircloak_ci_fix_permissions.service
[Unit]
Description=Permission fixer, see fix_permissions.sh for details

[Service]
Type=oneshot
ExecStart=/root/fix_permissions.sh
User=root
