# ci/production/aircloak_ci_restart_docker.service
[Unit]
Description=Restart docker service, see restart_docker.sh for details

[Service]
Type=oneshot
ExecStart=/root/restart_docker.sh
User=root
